# systemd/ziti-edge-tunnel-wsl.service
# OpenZiti tunnel — TUN mode for WSL2.
#
# WSL2 cannot use tproxy mode (no iptables mangle support in the default
# kernel). Instead we use ziti-edge-tunnel's TUN device (ziti0) with a
# companion ExecStartPost script that configures systemd-resolved split DNS.
#
# Install:
#   sudo cp systemd/ziti-edge-tunnel-wsl.service /etc/systemd/system/ziti-edge-tunnel.service
#   sudo cp systemd/ziti-resolved-setup.sh /usr/local/bin/ziti-resolved-setup.sh
#   sudo chmod +x /usr/local/bin/ziti-resolved-setup.sh
#   mkdir -p ~/.ziti/identities
#   # Copy your enrolled identity JSON to ~/.ziti/identities/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ziti-edge-tunnel
#
# Or use: scripts/install_wsl_tunnel.sh

[Unit]
Description=OpenZiti Edge Tunnel
Documentation=https://openziti.io
After=network-online.target systemd-resolved.service
Wants=network-online.target systemd-resolved.service

[Service]
Type=simple
ExecStart=/home/dev/bin/ziti-edge-tunnel run -i /home/dev/.ziti/identities/sean-wsl.json
ExecStartPost=/usr/local/bin/ziti-resolved-setup.sh
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
